# ubuntu/roles/ipset-deploy/files/ipset.service
[Unit]
Description=Loading IP Sets
Before=network.target
Before=netfilter-persistent.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/ipset -! -f /etc/ipset.conf restore
ExecReload=/sbin/ipset -! -f /etc/ipset.conf restore
ExecStop=/sbin/ipset flush

[Install]
WantedBy=multi-user.target
RequiredBy=netfilter-persistent.service
